# deployment/ybb-data-service.service
[Unit]
Description=YBB Data Management Service (Production)
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/var/www/ybb-data-service/current
Environment="PATH=/var/www/ybb-data-service/current/venv/bin"
Environment="FLASK_ENV=production"
Environment="PYTHONPATH=/var/www/ybb-data-service/current"
ExecStart=/var/www/ybb-data-service/current/venv/bin/gunicorn \
    --bind 0.0.0.0:5000 \
    --workers 4 \
    --worker-class sync \
    --timeout 300 \
    --log-level info \
    --access-logfile /var/www/ybb-data-service/current/logs/gunicorn_access.log \
    --error-logfile /var/www/ybb-data-service/current/logs/gunicorn_error.log \
    app:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target