# deployment/ybb-data-service-staging.service
[Unit]
Description=YBB Data Management Service (Staging)
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/var/www/ybb-data-service-staging/current
Environment="PATH=/var/www/ybb-data-service-staging/current/venv/bin"
Environment="FLASK_ENV=development"
Environment="PYTHONPATH=/var/www/ybb-data-service-staging/current"
ExecStart=/var/www/ybb-data-service-staging/current/venv/bin/gunicorn \
    --bind 0.0.0.0:5001 \
    --workers 2 \
    --worker-class sync \
    --timeout 300 \
    --log-level debug \
    --access-logfile /var/www/ybb-data-service-staging/current/logs/gunicorn_access.log \
    --error-logfile /var/www/ybb-data-service-staging/current/logs/gunicorn_error.log \
    app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target